# systemd/after-i8k.service
[Unit]
Description=Set fan speed to 0

[Service]
ExecStart=/usr/bin/i8kctl fan 0 0

[Install]
WantedBy=multi-user.target
